# etc/nsca-ng.socket
[Socket]
ListenStream=5668
BindIPv6Only=both

[Install]
WantedBy=sockets.target
